# consul-client2.conf
start on (filesystem and net-device-up IFACE=lo)
stop on runlevel [!2345]

exec /usr/bin/consul agent -data-dir /tmp/consul -bind 10.0.1.13 -join 10.0.1.10
